# raspberry-pi-setup/update.service
[Unit]
Description=ILSport dart-ad auto-update (git pull + restart)
After=network-online.target

[Service]
Type=oneshot
User=@USER@
EnvironmentFile=/etc/ilsport/env
ExecStart=@INSTALL_DIR@/raspberry-pi-setup/update.sh
